# debian/cron.daily
#!/bin/sh
set -e
if [ -x /usr/sbin/etckeeper ] && [ -e /etc/etckeeper/etckeeper.conf ]; then
	. /etc/etckeeper/etckeeper.conf
	if [ "$AVOID_DAILY_AUTOCOMMITS" != "1" ]; then
		if etckeeper unclean; then
			etckeeper commit "daily autocommit"
		fi
	fi
fi
